# infra/systemd/astra-caddy.service
[Unit]
Description=Astra Caddy reverse proxy
After=network-online.target astra-app@1.service astra-app@2.service
Requires=astra-app@1.service astra-app@2.service
Wants=network-online.target

[Service]
EnvironmentFile=/etc/astra/caddy.env
ExecStartPre=-/usr/bin/podman rm -f astra-caddy
ExecStart=/usr/bin/podman run --rm --name astra-caddy --network host -v /etc/astra/Caddyfile:/etc/caddy/Caddyfile:ro,Z -v /var/lib/caddy:/data:Z ${CADDY_IMAGE}
ExecStop=/usr/bin/podman stop -t 10 astra-caddy
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
